# etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

sh /etc/openvpn/vpn.sh

/etc/init.d/openvpn start

exit 0
